# dist/dpkg/DEBIAN/sysinfo-exporter.postinst
#!/bin/sh

set -ex

systemctl enable sysinfo-exporter
