# dist/dpkg/debian/sysinfo-exporter.service
[Unit]
Description=Sysinfo Exporter
After=network.target auditd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/sysinfo-exporter --config /etc/sysinfo-exporter/config.toml
KillMode=process
Restart=always

[Install]
WantedBy=default.target
Alias=sysinfo-exporter.service